# provision/etc/init.d/detect-cpu
#!/sbin/runscript

depend() {
    need localmount
    after bootmisc
}

start() {
    ebegin "Fetch cpu number"
    echo "MAKEOPTS=\"-j$(nproc)\"" > /etc/portage/cpu.conf
    eend $?
}
